# deploy/supervisord.conf
[supervisord]
nodaemon=true

[program:web]
command=apache2-foreground
autostart=true
autorestart=true
user=runtime
stdout_logfile=/dev/stdout
stderr_logfile=/dev/stderr